# systemd/doc-specialist.service
[Unit]
Description=OpenClaw Doc Specialist Service
ConditionPathExists=%h/openclaw-operator/agents/doc-specialist/src/service.ts
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw-operator/agents/doc-specialist
Environment=ALLOW_DIRECT_SERVICE=true
ExecStart=/usr/bin/env bash -lc 'exec node ../../node_modules/tsx/dist/cli.mjs src/service.ts'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
