# systemd/orchestrator.service
[Unit]
Description=OpenClaw Orchestrator Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw-operator/orchestrator
Environment=NODE_ENV=production
Environment=PORT=3312
Environment=PROMETHEUS_PORT=9100
Environment=ORCHESTRATOR_CONFIG=%h/openclaw-operator/orchestrator_config.json
Environment=OPERATOR_UI_DIST_DIR=%h/openclaw-operator/operator-s-console/dist
ExecStart=%h/openclaw-operator/orchestrator/node_modules/.bin/tsx src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=default.target
