# systemd/orchestrator.service
[Unit]
Description=OpenClaw Orchestrator Service
After=network-online.target openclaw-redis.service
Wants=openclaw-redis.service

[Service]
Type=simple
WorkingDirectory=%h/openclaw-operator/orchestrator
Environment=NODE_ENV=production
ExecStart=/usr/bin/env bash -lc 'exec node --env-file=.env dist/index.js'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=default.target
